# .crontab
0 */1 * * * ~/.config/transmission-daemon/scripts/clear-torrents.sh
