# .crontab
0 */2 * * * ~/.config/transmission-daemon/scripts/clear-torrents.sh
