# .crontab
0 */1 * * * ~/.config/transmission-daemon/scripts/clear-torrents.sh
0 */1 * * * ~/.config/conky/scripts/refresh-battery.sh
