# cron/.crontab
0 */3 * * * ~/.config/transmission-daemon/scripts/clear-torrents.sh
